# systemd/notify-failure@.service
[Unit]
Description=Notify that %i hit start limit

[Service]
Type=oneshot
ExecStart="$HOME"/bin/notify-failure.sh %i
