# systemd/notify-absent-from-steam.service
[Unit]
Description=Notify when absent from Steam

[Service]
Type=oneshot
ExecStart=/root/bin/notify-absent-from-steam.sh
