# systemd/tf2-dedicated@tf2-walkway.timer
[Unit]
Description=TF2 daily server restart (walkway)

[Timer]
OnCalendar=*-*-* 03:00:00 UTC
OnBootSec=1h
Persistent=true
Unit=tf2-dedicated@tf2-walkway.service

[Install]
WantedBy=timers.target
