# systemd/tf2-dedicated.timer
[Unit]
Description=TF2 daily server restart

[Timer]
OnCalendar=*-*-* 04:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
