# systemd/tf2-dedicated@tf2-jump.timer
[Unit]
Description=TF2 daily server restart (jump)

[Timer]
OnCalendar=*-*-* 05:30:00 Africa/Johannesburg
OnBootSec=2h20min
Persistent=true
Unit=tf2-dedicated@tf2-jump.service

[Install]
WantedBy=timers.target
