# systemd/tf2-dedicated@.service.d/failure-notify.conf
[Unit]
OnFailure=notify-failure@%n.service
